# .config/hypr/tools/random-wall/random-wall.timer
[Unit]
Description=Timer for random-wall service

[Timer]
OnUnitActiveSec=5min
OnbootSec=5min

[Install]
WantedBy=timers.target
