# .config/hypr/tools/random-wall/random-wall.service
[Unit]
Description=Script to randomize hyprpaper wallpapers

[Service]
Type=oneshot
ExecStart=/usr/bin/random-wall
